# overlay/etc/systemd/system/first-boot.service
[Unit]
Description=Run first boot setup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/first-boot-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target